# scripts/linux/health-backend.service.template
[Unit]
Description=Health System Flask backend
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
Group=__SERVICE_USER__
WorkingDirectory=__PROJECT_DIR__/backend
Environment=PYTHONUNBUFFERED=1
ExecStart=__PROJECT_DIR__/backend/.venv/bin/python -m waitress --listen=127.0.0.1:__BACKEND_PORT__ --threads=__BACKEND_THREADS__ wsgi:app
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
